# BTC_Prediction/deploy/systemd/sygnif-bybit-nl-feed.service.in
[Unit]
Description=Bybit market data to NeuroLinked (sygnif-swarm bundle)
After=network-online.target sygnif-neurolinked.service
Wants=network-online.target

[Service]
Type=simple
User=@@SERVICE_USER@@
Group=@@SERVICE_GROUP@@
WorkingDirectory=@@BUNDLE_ROOT@@/BTC_Prediction
EnvironmentFile=-@@BUNDLE_ROOT@@/BTC_Prediction/.env
EnvironmentFile=-@@BUNDLE_ROOT@@/BTC_Prediction/swarm_operator.env
Environment=SYGNIF_SWARM_BYBIT_MAINNET=1
Environment=BYBIT_NL_INCLUDE_SWARM=1
ExecStart=@@PYTHON@@ @@BUNDLE_ROOT@@/BTC_Prediction/scripts/bybit_nl_market_feed.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
